# ansible/roles/galaxy_exporter/templates/galaxy-exporter.service.j2
# {{ ansible_managed }}
[Unit]
Description=Galaxy Exporter — generate binary map data from Memgraph
After=network.target

[Service]
Type=oneshot
ExecStart={{ galaxy_exporter_binary_path }} \
  --memgraph-host={{ vpn_interface_address }} \
  --memgraph-port={{ memgraph_port }} \
  --output={{ galaxy_exporter_output_dir }} \
  --eddn-dsn=postgres://eddn_admin:{{ vault_eddn_db_password }}@{{ vpn_interface_address }}:5433/eddn_raw \
  --batch-size={{ galaxy_exporter_batch_size }} \
  --history-cycles={{ galaxy_exporter_history_cycles }} \
  --compress
StandardOutput=journal
StandardError=journal
SyslogIdentifier=galaxy-exporter
