# ansible/roles/galaxy_exporter/templates/galaxy-exporter.timer.j2
[Unit]
Description=Galaxy Static Data Exporter Timer
Documentation=man:systemd.timer(5)

[Timer]
# Run every 15 minutes
OnCalendar=*:0/15
# Randomize within 1 minute to avoid exact-minute load spikes
RandomizedDelaySec=60
# Run immediately on boot if last run was missed
Persistent=true
# Add some time for boot to settle
OnBootSec=2min

[Install]
WantedBy=timers.target
