# ansible/roles/galaxy_exporter/templates/galaxy-exporter.service.j2
[Unit]
Description=Galaxy Static Data Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ galaxy_exporter_bin_dir }}/galaxy-exporter \
    --output={{ galaxy_exporter_output_dir }} \
    --memgraph-host={{ memgraph_host | default('10.8.0.3') }} \
    --memgraph-port={{ memgraph_port | default(7687) }} \
    --batch-size={{ galaxy_exporter_batch_size | default(50000) }} \
    --compress={{ galaxy_exporter_compress | default('true') }} \
    --eddn-dsn=postgres://{{ galaxy_exporter_eddn_db_user }}:{{ galaxy_exporter_eddn_db_password }}@{{ galaxy_exporter_eddn_db_host }}:{{ galaxy_exporter_eddn_db_port }}/{{ galaxy_exporter_eddn_db_name }} \
    --history-cycles={{ galaxy_exporter_history_cycles }} \
    --history-min-cycle={{ galaxy_exporter_history_min_cycle | default(72) }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=galaxy-exporter

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ galaxy_exporter_output_dir }}

[Install]
WantedBy=multi-user.target
